# poker_bot/deploy/poker_bot.service
[Unit]
Description=Poker Bot + Mini App
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/poker_bot
ExecStart=/home/ubuntu/poker_bot/venv/bin/python main.py
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/poker_bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
